# Packaging/nanocharger.service
[Unit]
Description=NanoCharger USB charging scheduler
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/nanocharger daemon --config /etc/nanocharger/schedule.json --state-dir /var/lib/nanocharger
Restart=always
RestartSec=5s
# Runs as root by default. uhubctl requires root (or per-hub udev rules) to toggle VBUS, and a per-user
# Swift toolchain (e.g. installed via swiftly under $HOME) is otherwise unreadable by a system user.
# To run unprivileged: add `User=<name>`, set up uhubctl udev rules for that user, and ensure the Swift
# runtime libs are readable by that user (typically by mirroring them under /usr/local/lib/swift/linux).
StateDirectory=nanocharger
StateDirectoryMode=0755

[Install]
WantedBy=multi-user.target
